# systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper

[Service]
ExecStart=/home/sava/extra/wallswitcher.sh
