# systemd/user/wallpaper.timer
[Unit]
Description=Change wallpaper every 45 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Unit=wallpaper.service

[Install]
WantedBy=timers.target