# wallpapers/wallpaper.service
[Unit]
Description=Change wallpaper

[Service]
ExecStart=/home/sava/Extra/wallswitcher.sh
